# systemd/solen-backups@.service
[Unit]
Description=SOLEN backup profile %i

[Service]
Type=oneshot
Environment=SOLEN_PLAIN=1
EnvironmentFile=-%h/.config/solen/backups.env
EnvironmentFile=-/etc/solen/backups.env
ExecStart=%h/.local/bin/serverutils run backups/run -- run --profile %i --json
ExecStartPre=/usr/bin/env mkdir -p %h/.local/share/solen
StandardOutput=append:%h/.local/share/solen/backups-%i.ndjson
StandardError=journal
