# systemd/solen-update-check.timer
[Unit]
Description=Weekly SOLEN update check

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target

